# backend/launchd/com.qmd-index-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- managed-by: qmd-auto-context -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qmd-index-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@@HOME@@/.config/qmd/index_worker.sh</string>
    </array>
    <!-- 60초마다 dirty 큐 drain. 비었으면 즉시 종료. -->
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>QMD_DAEMON_PORT</key>
        <string>8483</string>
    </dict>
</dict>
</plist>
