# backend/launchd/com.qmd-logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- managed-by: qmd-auto-context -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qmd-logrotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@@HOME@@/.config/qmd/qmd-logrotate.sh</string>
    </array>

    <!-- 매일 04:30 에 로그 크기 가드. 임계(10MB) 미달이면 재시작 없이 즉시 종료. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>@@HOME@@</string>
    </dict>
</dict>
</plist>
